# 7.0/etc/supervisor/conf.d/20_openerp.conf
[program:openerp]
command=/usr/bin/openerp-server --config=/mnt/openerp-server.conf
user=openerp
